# debian/pyznap-send.timer
[Unit]
Description=Run pyznap send every 2 hours

[Timer]
OnCalendar=00/2:00
Persistent=true

[Install]
WantedBy=timers.target
